# app/config/anacrontab
* * * * * /<path to>/app/console kuma:nodes:cron --env=prod
